# run/systemd/generator/media-cdrom0.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
Where=/media/cdrom0
What=/dev/sr0
Type=udf,iso9660
Options=user,noauto
