# var/lib/dpkg/info/policycoreutils.conffiles
/etc/init.d/selinux-autorelabel
/etc/sestatus.conf
